# templates/docker-registry.service.erb
.include /lib/systemd/system/docker-registry.service
[Service]
User=<%= user %>
